# src/lkm/systemctl.service.in
[Unit]
Description=Load @PROJ_NAME@ loadable kernel module

[Service]
Type=oneshot
ExecStart=/sbin/modprobe @PROJ_NAME@
ExecStartPost=/bin/chown root:pruio /sys/devices/platform/@PROJ_NAME@/state
ExecStop=/sbin/modprobe -r @PROJ_NAME@
StandardOutput=tty
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
